# control/postinst
systemctl enable webgui.service
systemctl daemon-reload
systemctl start webgui.service