# control/prerm
systemctl stop webgui.service
systemctl disable webgui.service
systemctl daemon-reload
